# systemd/iprdump.service
[Unit]
Description=IBM Power Raid dump daemon

[Service]
Type=forking
ExecStart=/sbin/iprdump --daemon

[Install]
WantedBy=multi-user.target
