# systemd/iprdump.service.in
[Unit]
Description=IBM Power Raid dump daemon
Before=iprinit.service iprupdate.service

[Service]
Type=forking
ExecStart=@sbindir@/iprdump --daemon

[Install]
WantedBy=iprutils.target
