# systemd/iprinit.service.in
[Unit]
Description=IBM Power Raid init daemon

[Service]
Type=forking
ExecStart=@sbindir@/iprinit --daemon

[Install]
WantedBy=multi-user.target
