# systemd/iprupdate.service.in
[Unit]
Description=IBM Power Raid update daemon

[Service]
Type=forking
ExecStart=@sbindir@/iprupdate --daemon

[Install]
WantedBy=iprutils.target
